# systemd/cv-card-demo-kiosk.service.in
[Unit]
Description=CV Card Demo Local Kiosk Browser
After=graphical.target cv-card-demo-backend.service
Wants=cv-card-demo-backend.service
# Bind to the demo target so stop/restart of the target propagates here.
PartOf=cv-card-demo.target

[Service]
Type=simple
User=@USER@
Group=@GROUP@

Environment=DISPLAY=@DISPLAY@
Environment=XAUTHORITY=@XAUTHORITY@

ExecStart=@ROOT@/scripts/run_kiosk.sh

Restart=always
RestartSec=3

StandardOutput=journal
StandardError=journal

# Enabled via the target only — see backend unit for the rationale.
[Install]
WantedBy=cv-card-demo.target
