# systemd/cv-card-demo-kiosk.service.in
[Unit]
Description=CV Card Demo Local Kiosk Browser
After=graphical.target cv-card-demo-backend.service
Wants=cv-card-demo-backend.service

[Service]
Type=simple
User=@USER@
Group=@GROUP@

Environment=DISPLAY=@DISPLAY@
Environment=XAUTHORITY=@XAUTHORITY@

ExecStart=@ROOT@/scripts/run_kiosk.sh

Restart=always
RestartSec=3

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
