# systemd/cv-card-demo-backend.service.in
[Unit]
Description=CV Card Demo Backend and State Machine
Wants=network-online.target
After=network-online.target
# Bind this service to the demo target so `systemctl stop/restart
# cv-card-demo.target` propagates here (PartOf = stop/restart only, not start).
PartOf=cv-card-demo.target
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
User=@USER@
Group=@GROUP@
WorkingDirectory=@ROOT@

Environment=PYTHONUNBUFFERED=1
Environment=CV_CARD_DEMO_CONFIG=@ROOT@/config/config.yaml

ExecStart=@ROOT@/scripts/run_backend.sh

Restart=always
RestartSec=2

StandardOutput=journal
StandardError=journal

# Enabled via the target only (WantedBy=cv-card-demo.target), never directly
# into multi-user.target — otherwise `disable cv-card-demo.target` would not
# stop the service from autostarting on boot.
[Install]
WantedBy=cv-card-demo.target
